# data/playerctld.service
[Unit]
Description=Playerctl daemon
Documentation=man:playerctld(1)

[Service]
Type=exec
ExecStart=playerctld
